# etc/service/django/run
#!/bin/sh
set -e
cd $APP_ROOT
exec 2>&1
exec chpst -u www-data python3 -b $(which gunicorn) --config /etc/gunicorn/async.py wsgi
